# iems5722_app.conf
[program:iems5722_a3]
command = gunicorn -b localhost:8000 -w 4 iems5722_a3:app
directory = /home/iems5722
user = root
autostart = true
autorestart = true
stdout_logfile = /home/iems5722/app.log
redirect_stderr = true